# scripts/com.convergio.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.convergio.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/Roberdan/GitHub/convergio/daemon/target/release/convergio</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Roberdan/GitHub/convergio/daemon</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/Roberdan</string>
        <key>PATH</key>
        <string>/Users/Roberdan/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CONVERGIO_CLAUDE_BIN</key>
        <string>/Users/Roberdan/.local/bin/claude</string>
        <key>CONVERGIO_GH_BIN</key>
        <string>/opt/homebrew/bin/gh</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/convergio-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/convergio-daemon.err</string>
</dict>
</plist>
